# mynote.service
[Unit]
Description=MyNote2.0 Application
Documentation=https://github.com/Peter-cpu333/MyNote2.git
After=network.target
Wants=network.target

[Service]
Type=forking
User=mynote
Group=mynote
WorkingDirectory=/opt/mynote2.0
ExecStart=/opt/mynote2.0/start.sh
ExecStop=/opt/mynote2.0/stop.sh
ExecReload=/opt/mynote2.0/restart.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mynote2.0

# 环境变量
Environment=NODE_ENV=production
Environment=PYTHONPATH=/opt/mynote2.0/backend

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mynote2.0/logs /opt/mynote2.0/backend/database

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target